# wallpapers/.config/systemd/user/wallpaper@.timer
[Timer]
OnUnitInactiveSec=3h
AccuracySec=1s

[Unit]
Requires=wallpaper@.service

[Install]
WantedBy=timers.target
