# battery-monitor/.config/systemd/user/battery-monitor@.timer
[Timer]
OnUnitInactiveSec=2s
AccuracySec=1s

[Unit]
Requires=battery-monitor@.service

[Install]
WantedBy=timers.target
